# deploy/t2d-enable-https.service
[Unit]
Description=Enable HTTPS for T2D after DNS is ready
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/t2d-enable-https
