# deploy/t2d-enable-https.timer
[Unit]
Description=Check T2D DNS and enable HTTPS

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target
